# installers/macosx/xwhep.worker/installer/PckRoot/Library/LaunchDaemons/fr.in2p3.lal.xwhep.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
 <key>Label</key>
  <string>fr.in2p3.lal.xwhep.worker</string>
 <key>EnvironmentVariables</key>
  <dict>
   <key>JAVA_HOME</key>
   <string>/System/Library/Frameworks/JavaVM.framework/Versions/1.6/Home</string>
  </dict>
 <key>ProgramArguments</key>
  <array>
   <string>/System/Library/Frameworks/JavaVM.framework/Versions/1.6/Home/bin/java</string>
   <string>-Djavax.net.ssl.trustStore=/private/etc/xwhep.worker/keystore/xwhepworker.keys</string>
   <string>-Dxtremweb.cache=/tmp </string>
   <string>-Dxtremweb.cp=/private/etc/xwhep.worker/lib</string>
   <string>-Djava.awt.headless=true</string>
   <string>-cp</string>
   <string>/private/etc/xwhep.worker/lib/bcprov-jdk16-140.jar:/private/etc/xwhep.worker/lib/ibis-util-2.1.jar:/private/etc/xwhep.worker/lib/jetty-all-8.1.8.v20121106.jar:/private/etc/xwhep.worker/lib/jetty-util-8.1.8.v20121106.jar:/private/etc/xwhep.worker/lib/log4j-1.2.17.jar:/private/etc/xwhep.worker/lib/slf4j-api-1.7.2.jar:/private/etc/xwhep.worker/lib/slf4j-log4j12-1.7.2.jar:/private/etc/xwhep.worker/lib/smartsockets-1.4.jar:/private/etc/xwhep.worker/lib/xtremweb.jar</string>
   <string>xtremweb.common.HTTPLauncher</string>
   <string>--xwconfig</string>
   <string>/private/etc/xwhep.worker/conf/xtremweb.worker.conf</string>
  </array>
 <key>RunAtLoad</key>
  <true/>
 <key>Nice</key>
  <integer>10</integer>
 <key>UserName</key>
  <string>%SYSLOGIN%</string>
 <key>StandardOutPath</key>
  <string>/var/log/xwhep.worker.log</string>
 <key>StandardErrorPath</key>
  <string>/var/log/xwhep.worker.log</string>
 <key>WorkingDirectory</key>
  <string>/var/xwhep/worker/</string>
</dict>
</plist>
